# com.claude.orbital.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Service identity -->
    <key>Label</key>
    <string>com.claude.orbital</string>

    <!-- Binary and arguments -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/orbital/bin/orbital</string>
        <string>--base-dir</string>
        <string>/Users/michaelfinkler/orbital-data</string>
        <string>--config</string>
        <string>/opt/orbital/etc/orbital.env</string>
    </array>

    <!-- Run as orbital-guard user -->
    <key>UserName</key>
    <string>orbital-guard</string>
    <key>GroupName</key>
    <string>orbital</string>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/opt/orbital</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/opt/orbital</string>
        <key>PATH</key>
        <string>/opt/orbital/bin:/usr/local/bin:/usr/bin:/usr/sbin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Start at boot, restart on crash -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/opt/orbital/log/orbital.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/orbital/log/orbital.log</string>

    <!-- Resource limits -->
    <key>ProcessType</key>
    <string>Standard</string>
    <key>Nice</key>
    <integer>5</integer>

</dict>
</plist>
